# logstash/templates/supervisord.conf
[supervisord]
logfile=/dev/null;
pidfile=/var/run/supervisord.pid;
nodaemon=true

[program:healthcheck]
command=/app/bin/healthcheck
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
auto_start=true
autorestart=true

[program:confd]
command=confd -node %(ENV_ETCD)s -config-file /app/confd.toml -confdir /app
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
auto_start=true
autorestart=true

[program:logstash]
command=/opt/logstash/bin/logstash -f /opt/logstash/conf.d
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
auto_start=true
autorestart=true